# moirai-retention.timer
[Unit]
Description=Run Expected Path OpenSearch retention cleanup

[Timer]
OnCalendar=hourly
Persistent=true
Unit=moirai-retention.service

[Install]
WantedBy=timers.target
